# rootfs/etc/services.d/js_automations/run
#!/bin/sh
cd /app
exec node js_automations/server.js
